# crontab
# START CRON JOB
0 20 * * * python3 /app/script.py
# END CRON JOB